# bin/gel-reverse-tunnel.service
[Unit]
Description=GEL Controller reverse SSH tunnel
After=network-online.target
Wants=network-online.target
# Restart the tunnel indefinitely - don't give up after a burst of failures
StartLimitIntervalSec=0

[Service]
Type=simple
User=weedom
Group=weedom
EnvironmentFile=-/etc/gel-controller/tunnel.env
ExecStart=/opt/gel-controller/bin/gel-reverse-tunnel.sh
Restart=always
RestartSec=15
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
